# systemd/grub2-se-verifyserv.service
[Unit]


[Service]
ExecStart=/usr/local/sbin/grub2-se-verifyserv

[Install]
WantedBy=default.target
